# overlay-adv/lib/systemd/system/adv-sync-rtc-time.service
[Unit]
Description=Advantech sync rtc_time to system_time

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adv_timesync.sh

