# ubuntu18.04/overlay-adv/lib/systemd/system/adv-wifi-init.service
[Unit]
Description=advantech wifi init
After=networking.service

[Service]
Type=simple
ExecStart=/usr/bin/adv_wifi_init.sh

[Install]
Alias=adv-wifi-init.service
WantedBy=multi-user.target

